# recipes-utils/cg5317-utils/files/resetmodem.service
[Unit]

[Service]
ExecStart=/usr/bin/resetmodem.sh
User=root

[Install]
WantedBy=default.target

